# root/etc/services.d/ssh-mitm/run
#!/usr/bin/with-contenv bash

if [ ! -f /hostkey ]; then
    echo "creating new hostkey"
    /usr/bin/ssh-keygen -t rsa -f /hostkey -q -N ""
fi

exec ssh-mitm --host-key /hostkey ${SSH_MITM_ARGS[@]}
